# example_configs/example.cron.daily
#!/bin/bash -l
supervisorctl restart emdr_relay_go > /dev/null
