# example_configs/supervisord-relay.conf
[program:emdr_relay_go]
command=/home/emdr/emdr-relay-go/emdr-relay-go
numprocs=1
stopsignal=QUIT
autostart=true
autorestart=true
redirect_stderr=True
user=emdr
group=emdr